# etc/systemd/system/mariadb.service.d/unset_env_var_empty_fix.conf
[Service]
Environment="MYSQLD_OPTS="
Environment="_WSREP_NEW_CLUSTER="
